# config_systemd_user/backup.service
[Unit]
Description=Run daily backup

[Service]
ExecStart=/home/jk/.prog/backup.sh
